# createdb.conf
[program:createdb]
user=postgres
command=createdb -O docker gitlabhq_production
redirect_stderr=true
autorestart=false
priority=1
